# ha/extract_ha.ini
[program:extractor_ha]
command=/root/bimax_ha_migration/ha/extract_ha.sh
directory=/root/bimax_ha_migration/ha
autostart=true
autorestart=false
redirect_stderr=true
