# ha/monitor_ha.ini
[program:monitor_ha]
command=/root/bimax_ha_migration/ha/pm2_ha.sh monitor
directory=/root/bimax_ha_migration/ha
autostart=true
autorestart=false
redirect_stderr=true
